# deploy/systemd/overrides/klipper-host-mcu.service.d/override.conf
[Service]
# Allow Klipper host MCU to request realtime scheduling
LimitRTPRIO=99
LimitRTTIME=infinity
AmbientCapabilities=CAP_SYS_NICE
CapabilityBoundingSet=CAP_SYS_NICE
Nice=-10
